# packages/brain-server/debian/postinst
#!/bin/bash
# Post-installation script for brain-server
# This script runs after the package is installed

set -e

# Package information
PACKAGE_NAME="brain-server"
SERVICE_NAME="brain-server"
BINARY_PATH="/usr/local/bin/brain-server"
SERVICE_PATH="/etc/systemd/system/${SERVICE_NAME}.service"
CONFIG_DIR="/etc/brain-server"
DATA_DIR="/var/lib/brain-server"

# Colors for output (if terminal supports it)
if [ -t 1 ]; then
    RED='\033[0;31m'
    GREEN='\033[0;32m'
    YELLOW='\033[1;33m'
    NC='\033[0m'
else
    RED=''
    GREEN=''
    YELLOW=''
    NC=''
fi

log_info() {
    echo -e "${GREEN}[INFO]${NC} $*"
}

log_warn() {
    echo -e "${YELLOW}[WARN]${NC} $*"
}

log_error() {
    echo -e "${RED}[ERROR]${NC} $*"
}

# Create brain-server user if it doesn't exist
if ! id -u brain-server &>/dev/null; then
    log_info "Creating brain-server system user..."
    useradd --system --no-create-home --shell /bin/false brain-server || {
        log_warn "Could not create brain-server user (may already exist)"
    }
fi

# Create necessary directories
log_info "Creating directories..."
mkdir -p "${CONFIG_DIR}"
mkdir -p "${DATA_DIR}"
mkdir -p "${DATA_DIR}/db"
mkdir -p "${DATA_DIR}/models"
mkdir -p "${DATA_DIR}/logs"

# Set permissions
log_info "Setting permissions..."
chown -R brain-server:brain-server "${DATA_DIR}"
chmod 750 "${DATA_DIR}"
chmod 750 "${DATA_DIR}/db"
chmod 750 "${DATA_DIR}/models"
chmod 755 "${CONFIG_DIR}"

# Make binary executable
if [ -f "${BINARY_PATH}" ]; then
    chmod +x "${BINARY_PATH}"
    log_info "Binary permissions set"
else
    log_warn "Binary not found at ${BINARY_PATH}"
fi

# Install systemd service if it exists in the package
if [ -f "${SERVICE_PATH}" ]; then
    log_info "Installing systemd service..."

    # Reload systemd daemon
    systemctl daemon-reload || {
        log_warn "Could not reload systemd daemon"
    }

    # Enable the service
    systemctl enable ${SERVICE_NAME} || {
        log_warn "Could not enable ${SERVICE_NAME} service"
    }

    log_info "Systemd service installed and enabled"
else
    log_warn "Systemd service file not found at ${SERVICE_PATH}"
fi

# Create default config if it doesn't exist
CONFIG_FILE="${CONFIG_DIR}/config.toml"
if [ ! -f "${CONFIG_FILE}" ]; then
    log_info "Creating default configuration..."
    cat > "${CONFIG_FILE}" << 'EOF'
# Brain Server Configuration
# Generated by post-install script

[server]
host = "127.0.0.1"
port = 8765

[database]
path = "/var/lib/brain-server/db/brain.db"

[embedding]
model = "minishlab/potion-retrieval-32M"

[logging]
level = "info"
path = "/var/lib/brain-server/logs"
EOF
    chown brain-server:brain-server "${CONFIG_FILE}"
    chmod 640 "${CONFIG_FILE}"
    log_info "Default configuration created at ${CONFIG_FILE}"
fi

# Check if service should be started
# Don't start if this is an upgrade and the service is already running
if [ "$1" = "configure" ]; then
    if [ -z "$2" ]; then
        # First installation - ask if user wants to start
        log_info "First-time installation detected"
        if [ -t 0 ]; then
            # Interactive mode
            echo ""
            read -p "Would you like to start ${SERVICE_NAME} now? (Y/n) " -n 1 -r
            echo
            if [[ ! $REPLY =~ ^[Nn]$ ]]; then
                systemctl start ${SERVICE_NAME} || {
                    log_warn "Could not start ${SERVICE_NAME}"
                }
            fi
        else
            # Non-interactive mode - don't auto-start
            log_info "Run 'sudo systemctl start ${SERVICE_NAME}' to start the service"
        fi
    else
        # Upgrade - restart if service was running
        if systemctl is-active --quiet ${SERVICE_NAME}; then
            log_info "Restarting ${SERVICE_NAME} (service was running during upgrade)..."
            systemctl restart ${SERVICE_NAME} || {
                log_warn "Could not restart ${SERVICE_NAME}"
            }
        else
            log_info "Service ${SERVICE_NAME} was not running, not starting"
        fi
    fi
fi

# Summary
echo ""
log_info "═══════════════════════════════════════════════════════"
log_info "  Brain Server Installation Complete!"
log_info "═══════════════════════════════════════════════════════"
echo ""
log_info "Binary:     ${BINARY_PATH}"
log_info "Config:     ${CONFIG_FILE}"
log_info "Database:   ${DATA_DIR}/db/"
log_info "Logs:       ${DATA_DIR}/logs/"
echo ""
log_info "Commands:"
echo "  Start:   sudo systemctl start ${SERVICE_NAME}"
echo "  Stop:    sudo systemctl stop ${SERVICE_NAME}"
echo "  Status:  sudo systemctl status ${SERVICE_NAME}"
echo "  Logs:    sudo journalctl -u ${SERVICE_NAME} -f"
echo ""

exit 0
